# configs/etc/s6/sshd/run
#!/bin/bash

# setup users
/opt/setupusers.sh

# start sshd
/usr/sbin/sshd -D
