# configs/etc/s6/sshd/run
#!/bin/bash

# setup users, git and development environment
/opt/setupusers.sh
/opt/setupgit.sh
/opt/setupenv.sh

# start sshd
/usr/sbin/sshd -D
